# sentinel.service
[Unit]
Description=Sentinel v0.4 System Monitor Service
Documentation=https://github.com/VidGuiCode/sentinel
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
ExecStart=/usr/local/bin/sentinel --service
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Performance: Nice value for low priority
Nice=10

# Optional: Run as specific user
# User=sentinel
# Group=sentinel

[Install]
WantedBy=multi-user.target
